# agents/prometheus-vllm-sd-agent/deploy/aiunistor-prometheus-vllm-sd-agent.service
[Unit]
Description=AiUnistor Prometheus vLLM SD Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/aiunistor/prometheus-vllm-sd-agent
EnvironmentFile=/etc/aiunistor/prometheus-vllm-sd-agent.env
ExecStart=/opt/aiunistor/prometheus-vllm-sd-agent/.venv/bin/uvicorn app:app --host 0.0.0.0 --port 18083
Restart=always
RestartSec=3
KillSignal=SIGINT
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target